# splunk/rc.local
#!/bin/sh
/opt/splunk/bin/splunk start --accept-license --seed-passwd changeme
exit 0
