# splunk/rc.local
#!/bin/sh
/opt/splunk/bin/splunk start --accept-license --answer-yes
exit 0
